# deploy/tianquan-backend.service
[Unit]
Description=天权留学选校推荐引擎 (Backend)
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/tianquan/backend
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 3470 --workers 2
Restart=always
RestartSec=5
Environment=TIANQUAN_DB_PATH=/home/admin/.openclaw/workspace-study-abroad/study-abroad-advisor/data/advisor.db
Environment=TIANQUAN_DATA_DIR=/home/admin/.openclaw/workspace-study-abroad/study-abroad-advisor/data

[Install]
WantedBy=multi-user.target
